# store/works/Crupest.SecretTool/tools/crupest-secret-tool.service
[Unit]
Description=crupest v2ray service

[Service]
ExecStart=%h/.local/bin/Crupest.SecretTool

[Install]
WantedBy=default.target
